# remove_services.service
[Unit]
Description=Remove specified service files

[Service]
ExecStart=/usr/local/bin/remove_services.sh
